# projects/agent-dashboard/deploy/agent-dashboard.service
[Unit]
Description=Agent Dashboard Service
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/agent-dashboard
Environment=NODE_ENV=production
Environment=USE_SQLITE=true
ExecStart=/usr/bin/pm2 start /var/www/agent-dashboard/ecosystem.config.js --no-daemon
Restart=always

[Install]
WantedBy=multi-user.target
